# files/configs/crontab
00 00 * * * /bin/bash -c /home/steam/files/scripts/cleanlogs.sh